# deploy/supervisor/ziifra-worker.conf
; ZIIFRA queue worker — install to /etc/supervisor/conf.d/ziifra-worker.conf
; Requires QUEUE_CONNECTION=database (or redis) in .env

[program:ziifra-worker]
process_name=%(program_name)s_%(process_num)02d
command=php APP_ROOT/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=APP_ROOT/storage/logs/worker.log
stopwaitsecs=3600
